# zabbix-master/zabbix-speedtest-lan/speedtest-lan.cron
*/10 * * * * zabbix /etc/zabbix/script/speedtest-lan.sh -a
