# zabbix-master/zabbix-speedtest-lan/speedtest-lan.timer
[Unit]
Description=Run a LAN speedtest every 10 minutes

[Timer]
OnCalendar=*:0/10
# RandomizedDelaySec=30

[Install]
WantedBy=timers.target
