# zabbix-master/zabbix-speedtest-lan-sender/speedtest-lan.cron
*/10 * * * * zabbix /etc/zabbix/script/zabbix_speedtest-lan.sh #Send speedtest data to Zabbix